# home/dot-config/systemd/user/switch-gnome-wallpaper.service
[Unit]
Description=Switch GNOME wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/fish -c "switch-gnome-wallpaper"
