# home/dot-config/systemd/user/switch-gnome-wallpaper.timer
[Unit]
Description=Hourly switch of GNOME wallpaper

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
